# ansible-2017/playbooks/roles/videobox/templates/video-status/video-status.service
[Unit]
Description=video-status service
After=video-streamer.service

[Service]
ExecStart=/usr/local/bin/status.py
ExecReload=/bin/kill -HUP $MAINPID
KillMode=control-group
KillSignal=SIGKILL
Restart=always
WatchdogSec=3600

[Install]
WantedBy=multi-user.target
