# ansible/playbooks/roles/power_supply/templates/power_supply_metrics.timer
[Unit]
Description=power_supply metrics

[Timer]
OnUnitActiveSec=10s
OnBootSec=10s

[Install]
WantedBy=timers.target
