# ansible-2017/playbooks/roles/voctop/files/units/streamcast_backup.service
[Unit]
Description= backup rtmp stream cast service

[Service]
ExecStart=/opt/scripts/streamcast_backup.sh
ExecReload=/bin/kill -HUP $MAINPID
KillMode=control-group
KillSignal=SIGKILL
Restart=on-failure

[Install]
WantedBy=multi-user.target
