# ansible/playbooks/roles/video-monitoring/files/monitoring.cron
* * * * *	root	/usr/local/bin/do_checks.sh
