# ansible/playbooks/roles/video-box-mixer/files/audio-fetcher/audio-fetcher-1.service
[Unit]

Description=audio-fetcher chan 1

[Service]
ExecStart=/usr/local/bin/audio-fetcher 1
ExecReload=/bin/kill -HUP $MAINPID
KillMode=control-group
KillSignal=SIGKILL
Restart=on-failure

[Install]
WantedBy=multi-user.target
