# ansible/playbooks/roles/video-box/files/video-streamer/rc.local
#!/bin/sh

cpupower frequency-set -g performance -u 2000M -d 1900M
