# ansible/playbooks/roles/video-box-mixer/files/video-streamer/rc.local
#!/bin/sh

cpupower frequency-set -g performance -u 2900M -d 1900M
