# ansible-2017/playbooks/roles/stream-dumper/templates/streamdump.service
[Unit]
Description={{ item.host }} stream dump service

[Service]
ExecStart=/usr/local/bin/streamdump.sh {{ item.source }} /mnt/storage/{{ item.host }}/
ExecReload=/bin/kill -HUP $MAINPID
KillMode=control-group
KillSignal=SIGKILL
Restart=on-failure

[Install]
WantedBy=multi-user.target
