# ansible/playbooks/roles/video-stream-dump/templates/streamdump.target
[Unit]
Description=Streamdump target allowing to start/stop all streamdump-*.service instances at once
[Install]
WantedBy=multi-user.target
